# deploy/systemd/aureon-organism.service
# Aureon Organism daemon (connectome sweep/heal/weave → 100%) — standalone unit.
# Adjust WorkingDirectory / ExecStart / User to your install prefix.
[Unit]
Description=Aureon Organism daemon — connectome breath (sweep/heal/weave)
After=network-online.target
Wants=network-online.target
PartOf=aureon.target

[Service]
Type=simple
User=aureon
WorkingDirectory=/opt/aureon
Environment=PATH=/opt/aureon/.venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/opt/aureon/.env
ExecStart=/opt/aureon/.venv/bin/python -u -m aureon.core.organism_daemon
Restart=always
RestartSec=5

[Install]
WantedBy=aureon.target
