# production/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/aureon/logs/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root
childlogdir=/aureon/logs
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin",PYTHONUNBUFFERED="1",AUREON_ENABLE_AUTONOMOUS_CONTROL="1"
minfds=8192
minprocs=32768

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /aureon/config/supervisor.d/*.conf

# ═══════════════════════════════════════════════════════════════════════════
# CORE SYSTEMS - CRITICAL ORCHESTRATION (priority 1-9)
# ═══════════════════════════════════════════════════════════════════════════

[program:queen-eternal-machine]
command=python /aureon/app/queen_eternal_machine.py
autostart=true
autorestart=true
priority=1
stderr_logfile=/aureon/logs/queen-eternal-machine.err.log
stdout_logfile=/aureon/logs/queen-eternal-machine.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=15
stopwaitsecs=10

[program:trade-activation]
command=python /aureon/app/aureon_trade_activation.py
autostart=true
autorestart=true
priority=2
stderr_logfile=/aureon/logs/trade-activation.err.log
stdout_logfile=/aureon/logs/trade-activation.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=17
stopwaitsecs=10

[program:queen-consciousness]
command=python /aureon/app/aureon_queen_consciousness.py
autostart=true
autorestart=true
priority=3
stderr_logfile=/aureon/logs/queen-consciousness.err.log
stdout_logfile=/aureon/logs/queen-consciousness.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=19
stopwaitsecs=10

[program:orca-dual-hunter]
command=python /aureon/app/orca_dual_hunter.py
autostart=true
autorestart=true
priority=4
stderr_logfile=/aureon/logs/orca-dual-hunter.err.log
stdout_logfile=/aureon/logs/orca-dual-hunter.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=21
stopwaitsecs=10

[program:queen-power-redistribution]
command=python /aureon/app/queen_power_redistribution.py
autostart=true
autorestart=true
priority=5
stderr_logfile=/aureon/logs/queen-power-redistribution.err.log
stdout_logfile=/aureon/logs/queen-power-redistribution.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=23
stopwaitsecs=10

[program:labyrinth-snowball]
command=python /aureon/app/labyrinth_snowball_engine.py
autostart=true
autorestart=true
priority=6
stderr_logfile=/aureon/logs/labyrinth-snowball.err.log
stdout_logfile=/aureon/logs/labyrinth-snowball.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=25
stopwaitsecs=10

[program:aureon-51-live]
command=python /aureon/app/aureon_51_live.py
autostart=true
autorestart=true
priority=7
stderr_logfile=/aureon/logs/aureon-51-live.err.log
stdout_logfile=/aureon/logs/aureon-51-live.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=27
stopwaitsecs=10

[program:micro-profit-labyrinth]
command=python /aureon/app/micro_profit_labyrinth.py
autostart=true
autorestart=true
priority=8
stderr_logfile=/aureon/logs/micro-profit-labyrinth.err.log
stdout_logfile=/aureon/logs/micro-profit-labyrinth.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=29
stopwaitsecs=10

[program:kraken-ecosystem]
command=python /aureon/app/aureon_kraken_ecosystem.py
autostart=true
autorestart=true
priority=9
stderr_logfile=/aureon/logs/kraken-ecosystem.err.log
stdout_logfile=/aureon/logs/kraken-ecosystem.out.log
environment=PYTHONUNBUFFERED=1,AUREON_DRY_RUN=false,AUREON_ENABLE_AUTONOMOUS_CONTROL=1
stopasgroup=true
killasgroup=true
startsecs=31
stopwaitsecs=10

# ═══════════════════════════════════════════════════════════════════════════
# SUB-SYSTEMS - SPECIALIZED TRADERS (priority 10-37)
# ═══════════════════════════════════════════════════════════════════════════

[program:ira-sniper-mode]
command=python /aureon/app/ira_sniper_mode.py
autostart=true
autorestart=true
priority=10
stderr_logfile=/aureon/logs/ira-sniper-mode.err.log
stdout_logfile=/aureon/logs/ira-sniper-mode.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=33
stopwaitsecs=5

[program:multi-pair-live]
command=python /aureon/app/aureon_multi_pair_live.py
autostart=true
autorestart=true
priority=11
stderr_logfile=/aureon/logs/multi-pair-live.err.log
stdout_logfile=/aureon/logs/multi-pair-live.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=35
stopwaitsecs=5

[program:orca-smart-kill]
command=python /aureon/app/orca_smart_kill_cycle.py
autostart=true
autorestart=true
priority=12
stderr_logfile=/aureon/logs/orca-smart-kill.err.log
stdout_logfile=/aureon/logs/orca-smart-kill.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=37
stopwaitsecs=5

[program:queen-execute]
command=python /aureon/app/aureon_queen_execute.py
autostart=true
autorestart=true
priority=13
stderr_logfile=/aureon/logs/queen-execute.err.log
stdout_logfile=/aureon/logs/queen-execute.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=39
stopwaitsecs=5

[program:irish-patriots]
command=python /aureon/app/irish_patriot_scouts.py
autostart=true
autorestart=true
priority=14
stderr_logfile=/aureon/logs/irish-patriots.err.log
stdout_logfile=/aureon/logs/irish-patriots.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=41
stopwaitsecs=5

[program:queen-power-system]
command=python /aureon/app/queen_autonomous_power_system.py
autostart=true
autorestart=true
priority=15
stderr_logfile=/aureon/logs/queen-power-system.err.log
stdout_logfile=/aureon/logs/queen-power-system.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=43
stopwaitsecs=5

[program:queen-unified-startup]
command=python /aureon/app/aureon_queen_unified_startup.py
autostart=true
autorestart=true
priority=16
stderr_logfile=/aureon/logs/queen-unified-startup.err.log
stdout_logfile=/aureon/logs/queen-unified-startup.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=45
stopwaitsecs=5

[program:btc-trader]
command=python /aureon/app/aureon_btc_trader.py
autostart=true
autorestart=true
priority=17
stderr_logfile=/aureon/logs/btc-trader.err.log
stdout_logfile=/aureon/logs/btc-trader.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=47
stopwaitsecs=5

[program:orca-kill-executor]
command=python /aureon/app/orca_kill_executor.py
autostart=true
autorestart=true
priority=18
stderr_logfile=/aureon/logs/orca-kill-executor.err.log
stdout_logfile=/aureon/logs/orca-kill-executor.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=49
stopwaitsecs=5

[program:quantum-quackers]
command=python /aureon/app/quantum_quackers.py
autostart=true
autorestart=true
priority=19
stderr_logfile=/aureon/logs/quantum-quackers.err.log
stdout_logfile=/aureon/logs/quantum-quackers.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=51
stopwaitsecs=5

[program:aureon-qgita]
command=python /aureon/app/aureon_qgita.py
autostart=true
autorestart=true
priority=20
stderr_logfile=/aureon/logs/aureon-qgita.err.log
stdout_logfile=/aureon/logs/aureon-qgita.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=53
stopwaitsecs=5

[program:profit-mesh]
command=python /aureon/app/aureon_profit_mesh.py
autostart=true
autorestart=true
priority=21
stderr_logfile=/aureon/logs/profit-mesh.err.log
stdout_logfile=/aureon/logs/profit-mesh.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=55
stopwaitsecs=5

[program:celtic-warfare]
command=python /aureon/app/celtic_warfare_live.py
autostart=true
autorestart=true
priority=22
stderr_logfile=/aureon/logs/celtic-warfare.err.log
stdout_logfile=/aureon/logs/celtic-warfare.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=57
stopwaitsecs=5

[program:orca-quantum-stream]
command=python /aureon/app/orca_quantum_live_stream.py
autostart=true
autorestart=true
priority=23
stderr_logfile=/aureon/logs/orca-quantum-stream.err.log
stdout_logfile=/aureon/logs/orca-quantum-stream.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=59
stopwaitsecs=5

[program:queen-true-consciousness]
command=python /aureon/app/aureon_queen_true_consciousness.py
autostart=true
autorestart=true
priority=24
stderr_logfile=/aureon/logs/queen-true-consciousness.err.log
stdout_logfile=/aureon/logs/queen-true-consciousness.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=61
stopwaitsecs=5

[program:profit-now]
command=python /aureon/app/aureon_profit_now.py
autostart=true
autorestart=true
priority=25
stderr_logfile=/aureon/logs/profit-now.err.log
stdout_logfile=/aureon/logs/profit-now.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=63
stopwaitsecs=5

[program:safe-profit-trader]
command=python /aureon/app/aureon_safe_profit_trader.py
autostart=true
autorestart=true
priority=26
stderr_logfile=/aureon/logs/safe-profit-trader.err.log
stdout_logfile=/aureon/logs/safe-profit-trader.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=65
stopwaitsecs=5

[program:momentum-hunter]
command=python /aureon/app/aureon_live_momentum_hunter.py
autostart=true
autorestart=true
priority=27
stderr_logfile=/aureon/logs/momentum-hunter.err.log
stdout_logfile=/aureon/logs/momentum-hunter.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=67
stopwaitsecs=5

[program:turbo-snowball]
command=python /aureon/app/turbo_snowball_engine.py
autostart=true
autorestart=true
priority=28
stderr_logfile=/aureon/logs/turbo-snowball.err.log
stdout_logfile=/aureon/logs/turbo-snowball.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=69
stopwaitsecs=5

[program:queen-validated-trader]
command=python /aureon/app/aureon_queen_validated_trader.py
autostart=true
autorestart=true
priority=29
stderr_logfile=/aureon/logs/queen-validated-trader.err.log
stdout_logfile=/aureon/logs/queen-validated-trader.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=71
stopwaitsecs=5

[program:power-station-turbo]
command=python /aureon/app/aureon_power_station_turbo.py
autostart=true
autorestart=true
priority=30
stderr_logfile=/aureon/logs/power-station-turbo.err.log
stdout_logfile=/aureon/logs/power-station-turbo.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=73
stopwaitsecs=5

[program:quantum-warfare]
command=python /aureon/app/quantum_warfare_live.py
autostart=true
autorestart=true
priority=31
stderr_logfile=/aureon/logs/quantum-warfare.err.log
stdout_logfile=/aureon/logs/quantum-warfare.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=75
stopwaitsecs=5

[program:s5-live-trader]
command=python /aureon/app/s5_live_trader.py
autostart=true
autorestart=true
priority=32
stderr_logfile=/aureon/logs/s5-live-trader.err.log
stdout_logfile=/aureon/logs/s5-live-trader.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=77
stopwaitsecs=5

[program:live-conversion-trader]
command=python /aureon/app/live_conversion_trader.py
autostart=true
autorestart=true
priority=33
stderr_logfile=/aureon/logs/live-conversion-trader.err.log
stdout_logfile=/aureon/logs/live-conversion-trader.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=79
stopwaitsecs=5

[program:neural-revenue-orchestrator]
command=python /aureon/app/aureon_neural_revenue_orchestrator.py
autostart=true
autorestart=true
priority=34
stderr_logfile=/aureon/logs/neural-revenue-orchestrator.err.log
stdout_logfile=/aureon/logs/neural-revenue-orchestrator.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=81
stopwaitsecs=5

[program:momentum-snowball]
command=python /aureon/app/momentum_snowball_engine.py
autostart=true
autorestart=true
priority=35
stderr_logfile=/aureon/logs/momentum-snowball.err.log
stdout_logfile=/aureon/logs/momentum-snowball.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=83
stopwaitsecs=5

[program:aureon-mesh-live]
command=python /aureon/app/aureon_mesh_live.py
autostart=true
autorestart=true
priority=36
stderr_logfile=/aureon/logs/aureon-mesh-live.err.log
stdout_logfile=/aureon/logs/aureon-mesh-live.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=85
stopwaitsecs=5

[program:aureon-unified]
command=python /aureon/app/aureon_unified.py
autostart=true
autorestart=true
priority=37
stderr_logfile=/aureon/logs/aureon-unified.err.log
stdout_logfile=/aureon/logs/aureon-unified.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=87
stopwaitsecs=5

[program:aureon-infinite-kraken]
command=python /aureon/app/aureon_infinite_kraken.py
autostart=true
autorestart=true
priority=38
stderr_logfile=/aureon/logs/aureon-infinite-kraken.err.log
stdout_logfile=/aureon/logs/aureon-infinite-kraken.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=89
stopwaitsecs=5

[program:aureon-the-play]
command=python /aureon/app/aureon_the_play.py
autostart=true
autorestart=true
priority=39
stderr_logfile=/aureon/logs/aureon-the-play.err.log
stdout_logfile=/aureon/logs/aureon-the-play.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=91
stopwaitsecs=5

[program:snowball-conversion]
command=python /aureon/app/snowball_conversion_engine.py
autostart=true
autorestart=true
priority=40
stderr_logfile=/aureon/logs/snowball-conversion.err.log
stdout_logfile=/aureon/logs/snowball-conversion.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=93
stopwaitsecs=5

[program:compound-kelly-trader]
command=python /aureon/app/aureon_compound_kelly_trader.py
autostart=true
autorestart=true
priority=41
stderr_logfile=/aureon/logs/compound-kelly-trader.err.log
stdout_logfile=/aureon/logs/compound-kelly-trader.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=95
stopwaitsecs=5

[program:aureon-btc-v2]
command=python /aureon/app/aureon_btc_v2.py
autostart=true
autorestart=true
priority=42
stderr_logfile=/aureon/logs/aureon-btc-v2.err.log
stdout_logfile=/aureon/logs/aureon-btc-v2.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=97
stopwaitsecs=5

[program:aureon-the-play-old]
command=python /aureon/app/aureon_the_play_old.py
autostart=true
autorestart=true
priority=43
stderr_logfile=/aureon/logs/aureon-the-play-old.err.log
stdout_logfile=/aureon/logs/aureon-the-play-old.out.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=99
stopwaitsecs=5

# ═══════════════════════════════════════════════════════════════════════════
# WORKER GROUPS
# ═══════════════════════════════════════════════════════════════════════════

[group:core-systems]
programs=queen-eternal-machine,trade-activation,queen-consciousness,orca-dual-hunter,queen-power-redistribution,labyrinth-snowball,aureon-51-live,micro-profit-labyrinth,kraken-ecosystem
priority=1

[group:specialized-traders]
programs=ira-sniper-mode,multi-pair-live,orca-smart-kill,queen-execute,irish-patriots,queen-power-system,queen-unified-startup,btc-trader,orca-kill-executor,quantum-quackers,aureon-qgita,profit-mesh,celtic-warfare,orca-quantum-stream,queen-true-consciousness,profit-now,safe-profit-trader,momentum-hunter,turbo-snowball,queen-validated-trader,power-station-turbo,quantum-warfare,s5-live-trader,live-conversion-trader,neural-revenue-orchestrator,momentum-snowball,aureon-mesh-live,aureon-unified,aureon-infinite-kraken,aureon-the-play,snowball-conversion,compound-kelly-trader,aureon-btc-v2,aureon-the-play-old
priority=2

[group:all-trading-systems]
programs=queen-eternal-machine,trade-activation,queen-consciousness,orca-dual-hunter,queen-power-redistribution,labyrinth-snowball,aureon-51-live,micro-profit-labyrinth,kraken-ecosystem,ira-sniper-mode,multi-pair-live,orca-smart-kill,queen-execute,irish-patriots,queen-power-system,queen-unified-startup,btc-trader,orca-kill-executor,quantum-quackers,aureon-qgita,profit-mesh,celtic-warfare,orca-quantum-stream,queen-true-consciousness,profit-now,safe-profit-trader,momentum-hunter,turbo-snowball,queen-validated-trader,power-station-turbo,quantum-warfare,s5-live-trader,live-conversion-trader,neural-revenue-orchestrator,momentum-snowball,aureon-mesh-live,aureon-unified,aureon-infinite-kraken,aureon-the-play,snowball-conversion,compound-kelly-trader,aureon-btc-v2,aureon-the-play-old
priority=1
