# deploy/supervisord.linux.conf
; Aureon — full-system process definition for Linux (native, current paths).
; ─────────────────────────────────────────────────────────────────────────────
; Brings up the whole organism + trading swarm using the VERIFIED-CURRENT
; `python -m aureon.*` entry points (not the pre-refactor root-script names the
; older supervisord configs still reference). Launched by scripts/linux/aureon-up.sh
; or the deploy/systemd/aureon.service unit, both of which put the project venv's
; bin first on PATH (so `python` here is the venv interpreter) and set the run
; env (safe dry/paper by default; `--live` arms AUREON_LIVE_TRADING=1).
;
; SAFETY: live trading is OFF unless AUREON_LIVE_TRADING=1 (code default "0"); the
; hard boundary (conscience veto, AUREON_LOCAL_ACTIONS_ARMED / AUREON_SOUL_ACT stay
; unset) holds regardless — no autonomous irreversible execution.

[supervisord]
nodaemon=true
logfile=state/logs/supervisord.log
pidfile=state/supervisord.pid
childlogdir=state/logs
loglevel=info

[unix_http_server]
file=state/supervisor.sock

[supervisorctl]
serverurl=unix://state/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── organism core ────────────────────────────────────────────────────────────
[program:hnc-live-daemon]
command=python -u -m aureon.core.hnc_live_daemon
priority=10
autostart=true
autorestart=true
stdout_logfile=state/logs/hnc-live-daemon.log
redirect_stderr=true

[program:organism-daemon]
command=python -u -m aureon.core.organism_daemon
priority=11
autostart=true
autorestart=true
stdout_logfile=state/logs/organism-daemon.log
redirect_stderr=true

[program:operator]
command=python -u -m aureon.operator.operator_server
priority=12
environment=AUREON_OPERATOR_PORT="8790",AUREON_OPERATOR_HOST="0.0.0.0"
autostart=true
autorestart=true
stdout_logfile=state/logs/operator.log
redirect_stderr=true

; ── market + trading (DRY/PAPER by default — live only when AUREON_LIVE_TRADING=1) ──
[program:market-status-server]
command=python -u -m aureon.exchanges.unified_market_status_server --port 8791
priority=20
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/market-status-server.log
redirect_stderr=true

[program:ws-market-feeder]
command=python -u -m aureon.data_feeds.ws_market_data_feeder --binance --kraken --alpaca --capital --history-recording
priority=21
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/ws-market-feeder.log
redirect_stderr=true

[program:market-trader]
command=python -u -m aureon.exchanges.unified_market_trader --interval 30
priority=30
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/market-trader.log
redirect_stderr=true

[program:orca]
command=python -u -m aureon.bots.orca_complete_kill_cycle --autonomous --initial-capital 10000
priority=31
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/orca.log
redirect_stderr=true

[program:strategy-unity]
command=python -u -m aureon.trading.parallel_strategy_unity --watch
priority=32
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/strategy-unity.log
redirect_stderr=true

[program:queen-eternal-machine]
command=python -u -m aureon.queen.queen_eternal_machine --interval 10
priority=33
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/queen-eternal-machine.log
redirect_stderr=true

; ── autonomy / observer ────────────────────────────────────────────────────────
[program:mind-hub]
command=python -u -m aureon.autonomous.aureon_mind_thought_action_hub
priority=40
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/mind-hub.log
redirect_stderr=true

[program:self-questioning]
command=python -u -m aureon.autonomous.aureon_self_questioning_ai --loop
priority=41
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/self-questioning.log
redirect_stderr=true

[program:runtime-observer]
command=python -u -m aureon.autonomous.aureon_organism_runtime_observer --watch --refresh-core
priority=42
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/runtime-observer.log
redirect_stderr=true

[program:autonomous-self-run]
command=python -u -m aureon.autonomous.aureon_autonomous_self_run_loop --forever
priority=43
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/autonomous-self-run.log
redirect_stderr=true

; ── UIs ────────────────────────────────────────────────────────────────────────
[program:pro-dashboard]
command=python -u -m aureon.monitors.aureon_pro_dashboard
priority=50
autostart=%(ENV_AUREON_START_SWARM)s
autorestart=true
stdout_logfile=state/logs/pro-dashboard.log
redirect_stderr=true

[program:frontend]
command=npm run dev -- --host 0.0.0.0 --port 8081
directory=frontend
priority=60
autostart=%(ENV_AUREON_START_FRONTEND)s
autorestart=true
stdout_logfile=state/logs/frontend.log
redirect_stderr=true

; ── process groups (so `--organism-only` can start just the core) ───────────────
[group:organism]
programs=hnc-live-daemon,organism-daemon,operator

[group:swarm]
programs=market-status-server,ws-market-feeder,market-trader,orca,strategy-unity,queen-eternal-machine,mind-hub,self-questioning,runtime-observer,autonomous-self-run,pro-dashboard

[group:ui]
programs=frontend
